# config/systemd/ocean-one-engine.service
[Unit]
Description=Ocean ONE Daemon
After=network.target

[Service]
User=one
Type=simple
ExecStart=/home/one/bin/ocean.one -service engine
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
